# templates/default/strider.conf.erb
#!upstart

# This file was generated by Chef for <%= node['fqdn'] %>.
# Do NOT modify this file by hand! All changes WILL BE overwritten

description "Strider CD"
author "Chef for <%= node['fqdn'] %>"

start on startup
stop on shutdown

respawn
#respawn limit 5 30

setuid <%= @params['system-user'] %>
setgid <%= @params['system-user'] %>
env HOME="/home/<%= @params['system-user'] %>"
env DB_URI="<%= @params['db-uri'] %>"
env NODE_ENV=production
env SERVER_NAME="http://<%= @params['server-name'] %>"
# env HOST="<%= @params['server-name'] %>"
env PORT=<%= @params['server-port'] %>
env STRIDER_CLONE_DEST="/home/<%= @params['system-user'] %>/builds/"

env PLUGIN_GITHUB_APP_ID="<%= @params['github-appid'] %>"
env PLUGIN_GITHUB_APP_SECRET="<%= @params['github-appsecret'] %>"
# env PLUGIN_BITBUCKET_APP_KEY="<bitbucket app key>"
# env PLUGIN_BITBUCKET_APP_SECRET="<bitbuket secret>"
# env PLUGIN_BITBUCKET_HOSTNAME="http://<host>"

env SMTP_HOST="<%= @params['smtp-host'] %>"
env SMTP_PORT="<%= @params['smtp-port'] %>"
env SMTP_USER="<%= @params['smtp-user'] %>"
env SMTP_PASS="<%= @params['smtp-pass'] %>"
env SMTP_FROM="<%= @params['smtp-from'] %>"

script
    echo $$ > <%= @params['pid-file'] %>
    exec strider >> <%= @params['log-file'] %> 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> <%= @params['log-file'] %>
end script

pre-stop script
    rm <%= @params['pid-file'] %>
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> <%= @params['log-file'] %>
end script
